# addon/rootfs/etc/s6-overlay/s6-rc.d/amenitiz/finish
#!/command/execlineb
# Service finish script - called when service stops

multisubstitute
{
  importas -i ? EXIT_CODE
}

foreground
{
  echo "Domaine de Pipangaille Guest Manager stopped (exit code: ${EXIT_CODE})"
}

# Return the exit code
true
