# addon/rootfs/etc/s6-overlay/s6-rc.d/amenitiz/run
#!/command/execlineb -P
# Service run script for Home Assistant addon using s6-overlay

multisubstitute
{
  importas -i PORT PORT
}

foreground
{
  echo "Starting Domaine de Pipangaille Guest Manager on port ${PORT}"
}

cd /app
exec -c node src/server.js
